# profile.d/cluster.sh
export PATH=$PATH:/etc/cluster/scripts
export LC_ALL=en_US.UTF-8
